# lircd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/lircd --nodaemon --device=/dev/lirc0 --output=/run/lirc/lircd
